# vpn-tc-daemon.service
[Unit]
Description=OpenVPN TC Traffic Control Daemon
Documentation=man:tc(8)
After=network.target openvpn@server.service
Wants=openvpn@server.service
# 如果 OpenVPN 停止，这个服务也应该停止
BindsTo=openvpn@server.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/vpn-tc-daemon.sh
ExecReload=/bin/kill -HUP $MAINPID

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Watchdog 配置
WatchdogSec=30
# 脚本需要定期发送心跳，否则 systemd 会重启服务

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 安全加固
# ProtectSystem=strict
# ProtectHome=yes
# NoNewPrivileges=true
# PrivateTmp=yes

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vpn-tc-daemon

# 工作目录
WorkingDirectory=/var/log/openvpn

# 运行用户（必须是 root 才能操作 tc）
User=root
Group=root

# 环境变量
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# KillMode
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target